# secret-sluth.service
[Unit]
Description=Secret Sluth - Vault Secret Search Tool
Documentation=https://github.com/your-username/secret-sluth
After=network.target
Wants=network.target

[Service]
Type=simple
User=secret-sluth
Group=secret-sluth
WorkingDirectory=/opt/secret-sluth
Environment=PATH=/opt/secret-sluth/venv/bin
ExecStart=/opt/secret-sluth/venv/bin/python run.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/secret-sluth/logs /opt/secret-sluth/cache
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secret-sluth

# Environment variables
Environment=FLASK_ENV=production
Environment=FLASK_APP=app
Environment=PYTHONUNBUFFERED=1
Environment=SECRET_KEY=your-secret-key-change-in-production
Environment=VAULT_URL=
Environment=VAULT_TOKEN=

[Install]
WantedBy=multi-user.target
